# rootfs/etc/services.d/minidlna/run
#!/usr/bin/with-contenv sh

printf "starting minidlnad server...\m"
minidlnad -S -r -R -f /etc/minidlna/minidlna.conf
